# output/build/efl-1.20.6/systemd-services/ethumb.service
[Unit]
Description=Enlightenment Thumbnailer

[Service]
Type=dbus
BusName=org.enlightenment.Ethumb
ExecStart=/usr/bin/ethumbd
