# output/build/host-efl-1.20.6/systemd-services/ethumb.service
[Unit]
Description=Enlightenment Thumbnailer

[Service]
Type=dbus
BusName=org.enlightenment.Ethumb
ExecStart=/media/panpan/tools/buildroot-2018.02.8/output/host/bin/ethumbd
