# output/build/host-dbus-1.12.2/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/media/panpan/tools/buildroot-2018.02.8/output/host/var/run/dbus/system_bus_socket
